# scripts/pkg/run.loa.valetd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>run.loa.valetd</string>
    <key>Program</key>
    <string>/Applications/Valet.app/Contents/MacOS/valetd</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>__HOME__/Library/Logs/Valet/valetd.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>__HOME__/Library/Logs/Valet/valetd.stderr.log</string>
    <key>ProcessType</key>
    <string>Background</string>
    <key>AssociatedBundleIdentifiers</key>
    <string>run.loa.valet</string>
</dict>
</plist>
